# src/etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch --path.conf=/esvol/config
numprocs=1
autostart=true
autorestart=true